# config/supervisor_h4/bin_code.conf
[program:bin_code]
command=/var/www/mostbet/current/bin/console rabbitmq:consumer bin_code -m 1000 --without-signals --env=prod
process_name=%(process_num)s
numprocs=2
user=mostbet
autorestart=true
startretries=20
redirect_stderr=true
stdout_logfile=/var/log/project_log/bin_code.log
